# deploy/sub2api-datamanagementd.service
# Deprecated tombstone. datamanagementd has been removed from current Sub2API releases.
# Do not copy this unit to production. It intentionally fails to prevent
# accidental resurrection of the obsolete /tmp/sub2api-datamanagement.sock flow.

[Unit]
Description=Deprecated Sub2API Data Management Daemon

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo "datamanagementd has been removed from current Sub2API releases; do not install this service." >&2; exit 1'

[Install]
WantedBy=multi-user.target
